# supervisord.conf
[supervisord]
nodaemon=true

[program:kewatcher]
command=/bin/bash -c "kewatcher -c ./config/resque.yml"

[program:puma]
command=/bin/bash -c "bundle exec puma -C ./puma.rb"
