# supervisord.conf
[supervisord]
nodaemon=true

[program:kewatcher]
command=/bin/bash -c "kewatcher -c ./config/resque.yml"
